# cookbooks/stackful-node/templates/default/upstart/node-web.conf.erb
description "node-web server"

start on runlevel [2345]
respawn

setuid <%= node['stackful-node']['user'] %>
setgid stackful

env PYTHONPATH=<%= node["stackful-git"]["deployer-home"] %>
chdir <%= node['stackful-node']['app-home'] %>
exec <%= node['stackful-git']['deployer-home'] %>/bin/stackful-node-web /etc/stackful/node.json
